# DEBIAN/postinst
#!/bin/bash
set -e
if [ "$1" == "configure" ]; then
	systemctl enable system-installer
fi
