# DEBIAN/prerm
#!/bin/bash
set -e
set -o pipefail
#remove files generated by post-inst script
systemctl disable system-installer
